# alpine-influxdb/root/etc/services.d/influxdb/run
#!/usr/bin/with-contenv sh

exec s6-setuidgid influxdb influxd -config /etc/influxdb/influxdb.conf 2>$INFLUXDB_LOG_FILE;
