# alpine-php-debug/root/etc/services.d/php-fpm/run
#!/usr/bin/with-contenv sh

# start php-fpm;
php-fpm;
