# examples/user-consul-nginx-nodejs/root/etc/services.d/app/run
#!/usr/bin/with-contenv sh

cd /app

exec s6-setuidgid nodejs node server.js;
