# alpine-consul/root/etc/services.d/consul/finish
#!/usr/bin/with-contenv sh

s6-svscanctl -t /var/run/s6/services
