# alpine-rabbitmq/root/etc/services.d/rabbitmq/run
#!/usr/bin/with-contenv sh

# switch to the rabbitmq user and run the server
su rabbitmq -c rabbitmq-server;
