# alpine-rabbitmq/root/etc/services.d/rabbitmq/run
#!/usr/bin/with-contenv sh

./set-user-account

ulimit -n 1024
exec rabbitmq-server
