# examples/user-nginx-nodejs/root/etc/services.d/app/run
#!/usr/bin/with-contenv sh

cd /app

exec node server.js;
