# alpine-nginx/root/etc/services.d/nginx-access-log/run
#!/usr/bin/execlineb -P

redirfd -rnb 0 /var/run/s6/nginx-access-log-fifo
s6-envuidgid -D 32768:32768 nobody
s6-applyuidgid -U
exec -c
s6-log -bp -- t /var/log/nginx-access-logs
