# alpine-consul-redis/root/etc/services.d/redis/run
#!/usr/bin/with-contenv sh

s6-setuidgid redis redis-server /etc/redis.conf
