# alpine-grafana/root/etc/services.d/grafana/run
#!/usr/bin/with-contenv sh

exec s6-setuidgid grafana grafana-server --homepath /usr/share/grafana --config /etc/grafana/grafana.ini
