# examples/user-consul-agent/root/etc/services.d/consul/run
#!/usr/bin/with-contenv sh

# Bind to the external (LAN) IP
JOIN=`getent hosts consul-agent | awk '{print $1}'`

exec s6-setuidgid consul consul agent -config-dir /etc/consul.d -join ${JOIN};
