# examples/complete/app/root/etc/services.d/app/run
#!/usr/bin/with-contenv sh

cd /app

node index.js
