# alpine-apache/root/etc/services.d/apache/run
#!/usr/bin/with-contenv sh

mkdir -p /run/apache2

exec /usr/sbin/apachectl -DFOREGROUND;
